# xiboplayer/xiboplayer-chromium.service
[Unit]
Description=Xibo Digital Signage Player (Chromium Kiosk)
Documentation=https://github.com/linuxnow/xiboplayer-chromium
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/xiboplayer-chromium
Restart=always
RestartSec=5
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xiboplayer-chromium

[Install]
WantedBy=graphical-session.target
